# usr/local/etc/profile.d/bash_completion.sh
../../Cellar/bash-completion/1.3_3/etc/profile.d/bash_completion.sh